# service/tomcat/run
#!/bin/bash
set -e

exec catalina.sh run
